# ansible/files/webapp.socket
[Unit]
Description=Galaxy labs engine socket
[Socket]
ListenStream=/run/labs_engine.sock
[Install]
WantedBy=sockets.target
